# etc/systemd/system/docker.service.d/20-http-proxy.conf
[Service]
Environment="HTTP_PROXY=http://173.38.209.13:80/" "HTTPS_PROXY=http://173.38.209.13:80/" "NO_PROXY=192.168.7.193"
