# etc/systemd/system/docker.service.d/30-custom.conf
[Service]
Environment="DOCKER_OPTS=--cluster-advertise enp6s0:2375 --cluster-store etcd://192.168.7.193:2379"

